# homelab/Tools/automation/homelab-autoupdate.timer
[Unit]
Description=Run homelab-autoupdate daily at 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
