# homelab/noblenumbat-srv/yams/bb-kavita-sync.timer
[Unit]
Description=Run bb-kavita-sync hourly while Batman Beyond backlog downloads

[Timer]
OnBootSec=10min
OnUnitInactiveSec=1h
AccuracySec=5min

[Install]
WantedBy=timers.target
