# homelab/noblenumbat-srv/yams/media-import.timer
[Unit]
Description=Run media-import every 2 minutes

[Timer]
OnBootSec=60
OnUnitInactiveSec=2min
AccuracySec=30s

[Install]
WantedBy=timers.target
